# ansible/playbooks/roles/vhosts/grafana/templates/grafana-dash-pull.service.j2
[Unit]
Description=git pull dashboards

[Service]
Type=oneshot
WorkingDirectory={{ grafana_root_dir }}
ExecStart=/usr/bin/git pull --ff-only
