# ansible/playbooks/roles/vhosts/grafana/templates/grafana-dash-pull.timer.j2
[Unit]
Description=git pull dashboards every 5m

[Timer]
OnBootSec=30s
OnUnitActiveSec=5m
AccuracySec=30s

[Install]
WantedBy=timers.target
