# ansible/playbooks/roles/vhosts/node_exporter/templates/node-exporter.service
[Unit]
Description=Prometheus Node Exporter (tiny, LAN)
After=network.target

[Service]
User=nobody
Group=nogroup
ExecStart=/opt/metrics-agent/node_exporter \
  --web.listen-address=0.0.0.0:9100 \
  --collector.disable-defaults \
  --collector.cpu \
  --collector.meminfo \
  --collector.loadavg \
  --collector.filesystem \
  --collector.diskstats \
  --collector.netdev \
  --collector.netclass \
  --collector.uname \
  --collector.textfile \
  --collector.textfile.directory=/var/lib/node_exporter
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
